# .config/systemd/user/tasks.service
[Unit]
Description=Start tasks daemon

[Service]
ExecStart=%h/.local/bin/tasks
Nice=-20

[Install]
WantedBy=default.target
