# .config/autostart/restart-kmonad.desktop
[Desktop Entry]
Name=Restart Kmonad
Exec=rk
Terminal=false
Type=Application
NoDisplay=true
