# .config/autostart/restart-kmonad.desktop
[Desktop Entry]
Exec=rk
Name=Restart Kmonad
Type=Application
